# roles/archlinux/backup/templates/backup.timer.j2
[Unit]
Description=Scheduled backup

[Timer]
OnCalendar={{ backup_day }} {{ backup_time }}
AccuracySec=1h
Persistent=true

[Install]
WantedBy=timers.target
