# roles/archlinux/backup/templates/backup.service.j2
[Unit]
Description=Backup Service

[Service]
Type=simple
ExecStart={{ backup_script }}
User=root
Group=root
